# deploy/systemd/paperclip-recovery-monitor-backup.service
[Unit]
Description=PaperClip Recovery Monitor Backup (watches paperclip-recovery-monitor workflow liveness)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=python3 /home/sirrus/projects/BTC-Trade-Engine-PaperClip/scripts/paperclip_recovery_monitor_backup.py --json-summary
StandardOutput=journal
StandardError=journal
Environment=HOME=/home/sirrus
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/snap/bin:/home/sirrus/.local/bin
WorkingDirectory=/home/sirrus/projects/BTC-Trade-Engine-PaperClip
# Do NOT use EnvironmentFile=.env here: systemd's EnvironmentFile= overrides Environment=
# (empirically verified — opposite of older docs), and .env ships placeholder values
# (api.paperclip.example.com / dev_key / dev-company) that break local API access.
# The Python script calls load_dotenv(REPO_ROOT/".env") with override=False, so non-conflicting
# vars (POSTGRES_*, GH_TOKEN, …) still load while the systemd Environment= below wins.
Environment=PAPERCLIP_LISTEN_HOST=127.0.0.1
Environment=PAPERCLIP_LISTEN_PORT=3100
Environment=PAPERCLIP_API_KEY=local_trusted
Environment=PAPERCLIP_COMPANY_ID=73419cf3-bd37-4a7c-8782-311ccb47fced

[Install]
WantedBy=default.target
