# deploy/systemd/paperclip-recovery-monitor-backup.timer
[Unit]
Description=PaperClip Recovery Monitor Backup Timer (every 30 min, matching GH Actions :22/:52)
Requires=paperclip-recovery-monitor-backup.service

[Timer]
OnCalendar=*-*-* *:22:00
OnCalendar=*-*-* *:52:00
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
