# deploy/systemd/paperclip-impact-gate-scan-done.service
[Unit]
Description=Paperclip Impact Gate Scan-Done (5-min polling daemon — full gate)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/sirrus/projects/BTC-Trade-Engine-PaperClip/scripts/run_impact_gate_polling_daemon.sh --initial-scan
StandardOutput=journal
StandardError=journal
StandardInput=null
Environment=HOME=/home/sirrus
Environment=PATH=/home/sirrus/miniconda3/bin:/usr/local/bin:/usr/bin:/bin:/home/sirrus/.local/bin
EnvironmentFile=/home/sirrus/projects/BTC-Trade-Engine-PaperClip/.env
WorkingDirectory=/home/sirrus/projects/BTC-Trade-Engine-PaperClip

[Install]
WantedBy=default.target
