# scripts/backup-routine.service
[Unit]
Description=4-hourly database backup routine
Documentation=https://github.com/anomalyco/opencode
Wants=backup-routine.timer

[Service]
Type=oneshot
ExecStart=%h/projects/BTC-Trade-Engine-PaperClip/scripts/run_backup_routine.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
