# scripts/backup-routine.timer
[Unit]
Description=Backup routine timer (every 4 hours)
Documentation=https://github.com/anomalyco/opencode

[Timer]
OnCalendar=*:0/4
Persistent=true

[Install]
WantedBy=timers.target
