# scripts/opencode-watchdog.timer
[Unit]
Description=OpenCode watchdog timer (every 15 minutes)
Documentation=https://github.com/anomalyco/opencode

[Timer]
OnCalendar=*:0/15
Persistent=false

[Install]
WantedBy=timers.target
